# courier-authlib/courier-authlib.service.in
[Unit]
Description=Courier Authentication Library

[Service]
Type=forking
RemainAfterExit=true
ExecStart=@libexecdir@/courier-authlib/courier-authlib.sysvinit start
ExecStop=@libexecdir@/courier-authlib/courier-authlib.sysvinit stop

[Install]
WantedBy=multi-user.target
