# courier/courier.service.in
[Unit]
Description=Courier Mail Server
Wants=network-online.target
Wants=courier-authlib.service
After=network-online.target
After=courier-authlib.service

[Install]
WantedBy=mail-transfer-agent.target
WantedBy=multi-user.target

# First start could take several minutes, to generate the DH param file
[Service]
Type=forking
RemainAfterExit=true
ExecStart=@datadir@/courier.sysvinit start
ExecStop=@datadir@/courier.sysvinit stop
